# engine_cpp/deploy/ebbackup.timer
[Unit]
Description=Run ebbackup queue drain hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
